# configs/systemd/tunnelpanel.service
[Unit]
Description=TunnelPanel Server
After=network.target
Wants=network-online.target

[Service]
Type=simple
User=root
Group=root
ExecStart=/usr/local/bin/tunnelpanel
Restart=always
RestartSec=5
StandardOutput=journal
StandardError=journal
SyslogIdentifier=tunnelpanel

# Security hardening
NoNewPrivileges=false
ProtectSystem=false
PrivateTmp=true

# Environment
Environment=GIN_MODE=release

[Install]
WantedBy=multi-user.target
